# live-data.service
[Unit]
Description=Trading Live Data Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trading-data-service
Environment=PATH=/home/ubuntu/trading-data-service/venv/bin
ExecStart=/home/ubuntu/trading-data-service/venv/bin/python price_capture/live_data_launcher.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target